# launchd/com.timbre.morningcall.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.timbre.morningcall</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/node3/projects/voice_fun/scripts/morning_call.sh</string>
    </array>

    <!-- Fire every day at 07:00 local time (PST/PDT — launchd tracks DST automatically). -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Capture launchd-level stdout/stderr (the wrapper also tees into morning_call.log). -->
    <key>StandardOutPath</key>
    <string>/Users/node3/projects/voice_fun/logs/launchd_morning_call.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/node3/projects/voice_fun/logs/launchd_morning_call.log</string>

    <!-- CRITICAL: RunAtLoad is explicitly false — never fire on load. -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
